# templates/llmster.service.j2
[Unit]
Description=LM Studio headless API server (llmster)
After=network.target

[Service]
Type=simple
User={{ hermes_user }}
Environment=HOME={{ hermes_home }}
Environment=PATH={{ lmstudio_home }}/bin:{{ hermes_home }}/.local/bin:/usr/local/bin:/usr/bin:/bin
{% if hermes_model_api_key | default('') | trim | length > 0 %}
Environment=LM_API_TOKEN={{ hermes_model_api_key }}
{% endif %}
ExecStartPre=-{{ lmstudio_lms_cmd }} daemon up
ExecStart={{ lmstudio_lms_cmd }} server start --port {{ lmstudio_server_port | default(1234) }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
